# systemd/logi@.service
[Unit]
Description=Logitech R800 on %I

[Service]
Environment="args="
ExecStart=/usr/bin/logi $args %I
StandardOutput=journal
StandardError=journal
