# systemd/logi@.service.in
[Unit]
Description=Logitech R800 control on %I

[Service]
ExecStart=/usr/bin/@PROJECT_NAME@ %I
