# systemd/logi@.service
[Unit]
Description=Logitech R800 control on %I

[Service]
ExecStart=/usr/bin/logi %I
StandardOutput=journal
StandardError=journal
