# udev/50-logi.rules.in
ACTION=="add", SUBSYSTEM=="input", ATTRS{idVendor}=="046d", ATTRS{idProduct}=="c538", RUN+="@CMAKE_INSTALL_PREFIX@/bin/@PROJECT_NAME@ $env{DEVNAME}"
